# files/profile.d/bash_envar.sh
#
# /etc/profile.d/bash_envar.sh
# All environment vars that should always be made available

KDE_PREFIX=/opt/kde
KF5_PREFIX=/opt/kf5
LXQT_PREFIX=/opt/lxqt
LXQTDIR=/opt/lxqt
RAZORQT_PREFIX=/opt/razorqt
RAZORQTDIR=/opt/razorqt
EDE_PREFIX=/opt/enlightenment
EDEDIR=/opt/enlightenment
QT4DIR=/opt/qt4
QT5DIR=/opt/qt5

# The default, and typical, value for this is just the home directory
SYNCPATH=${HOME}

