# files/profile.d/bash_alias.sh
#!/bin/bash
#
# Begin /etc/profile.d/bash_alias.sh
# Sets the system-wide alias environment variables

alias ls='ls --color=auto'
alias grep='grep --color=auto'

# Source scripts
alias setqt4='source setqt4'
alias setqt5='source setqt5'
alias loadqt4='source loadqt4'
alias loadqt5='source loadqt5'
alias setqt4='source setkde'
alias setqt5='source setkf5'
alias setqt5='source setlxqt'

# End /etc/profile.d/bash_alias.sh
