# files/profile.d/pkgconfig.sh
#!/bin/bash
# Begin /etc/profile.d/pkgconfig.sh

#  PKG_CONFIG_PATH environment variable, and
#+ potentially anything else related to pkgconfig
export PKG_CONFIG_PATH="/usr/lib/pkgconfig"

# End /etc/profile.d/pkgconfig.sh
